# docker/crontab
# Run the script every 10 minutes
* * * * * /find-home >> /var/log/cron.log 2>&1
